# dump_wiki.cron
13 12 * * sat /var/www/www.ros.org/dump_wiki > /dev/null 2> /dev/null
