# dump_wiki.cron
13 12 * * 6 /var/www/www.ros.org/cron/dump_wiki > /dev/null 2> /dev/null
